# ansible/deploy_django/templates/systemctl/app.target
[Unit]
After=network.target
Wants={{ app_name }}_django.service {{ app_name }}_celery.service

[Install]
WantedBy=multi-user.target
